# deploy/holafresca-testing.service
[Unit]
Description=HolaFresca (testing branch, LAN)
After=network-online.target

[Service]
WorkingDirectory=%h/Documents/Programming/AI/HolaFresca-testing
# Land on the current integration head before starting. `-` so a sync hiccup
# (e.g. main repo mid-rebuild) never blocks the server from coming up.
ExecStartPre=-%h/Documents/Programming/AI/HolaFresca-testing/deploy/sync-integration.sh
ExecStart=%h/Documents/Programming/AI/HolaFresca-testing/.venv/bin/python run.py
Environment=HOST=0.0.0.0
Environment=PORT=8100
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
