# deploy/holafresca-testing-sync.service
[Unit]
Description=Sync HolaFresca testing checkout to integration head
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/Documents/Programming/AI/HolaFresca-testing/deploy/sync-integration.sh --restart
